# root/etc/s6-overlay/s6-rc.d/svc-avahi/run
#!/usr/bin/with-contenv bash

if [ "$START_AVAHI" = 'true' ]; then
  exec s6-notifyoncheck avahi-daemon --no-chroot
fi
