# root/etc/s6-overlay/s6-rc.d/svc-dbus/run
#!/usr/bin/with-contenv bash

exec \
  s6-setuidgid messagebus /usr/bin/dbus-daemon --system
