# root/etc/s6-overlay/s6-rc.d/svc-dbus/run
#!/usr/bin/with-contenv bash

echo "Starting dbus"
mkdir -p /var/run/dbus
exec s6-notifyoncheck dbus-daemon --system --nofork --nopidfile
